# dot/systemd/update_dns@.timer
[Unit]
Description=Scrape IP address, update DNS record on boot

[Timer]
OnBootSec=100
Unit=update_dns@%i.service

[Install]
WantedBy=timers.target
