# systemd/email/status-email-drbean@.service
[Unit]
Description=status email for %I to drbean

[Service]
Type=oneshot
ExecStart=/usr/local/bin/unit-error-email.sh "drbean@localhost" %i
User=nobody
Group=systemd-journal

